# etc/tilde.conf
# tilde - home management

description "tilde home manager"

start on runlevel [2345]
stop on runlevel [!2345]

respawn

exec /usr/bin/tilded -i60

